# infra/systemd/nullmaps-backup.timer
[Unit]
Description=Run NullMaps backup mid-week (Wed 03:30 UTC)

[Timer]
OnCalendar=Wed *-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
