# infra/systemd/nullmaps-refresh.service
[Unit]
Description=NullMaps weekly data refresh (OSM extract -> index/graph rebuild + backup)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/nullmaps
EnvironmentFile=-/opt/nullmaps/.env
ExecStart=/usr/bin/env bash /opt/nullmaps/infra/refresh.sh
